# suertech_agent.service
[Unit]
Description=Suertech AI Developer Agent
After=network.target

[Service]
# Бот будет работать от имени этого пользователя
User=ai_agent
Group=ai_agent

# Укажи путь, куда ты склонировал репозиторий
WorkingDirectory=/opt/my_agent

# Укажи путь к питону (рекомендуется использовать venv: /opt/my_agent/venv/bin/python)
ExecStart=/usr/bin/python3 -m src.bot.main

# Перезапуск при падении
Restart=always
RestartSec=10

# Перенаправление логов в syslog
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=suertech_agent

[Install]
WantedBy=multi-user.target
